# config/systemd/user/hugo.service
[Unit]
Description=hugo
StartLimitIntervalSec=0

[Service]
Restart=always
RestartSec=1s
WorkingDirectory=/home/michael/hugo
# Force native Go resolver
Environment=GODEBUG=netdns=go
ExecStart=/home/michael/go/bin/hugo serve -D --bind midna.monkey-turtle.ts.net --baseURL http://midna.monkey-turtle.ts.net:1313/

[Install]
WantedBy=default.target
